# crontab
0,20,40 * * * * dscacheutil -flushcache
0 0 * * * defaults write com.apple.dock persistent-apps -array
0 0 * * * /Users/macbook/Scripts/manage_packages.sh
0 0 * * 6 tmutil startbackup
